# container/php-fpm-worker.conf
[program:php5-fpm]
command = php-fpm
user = nginx 
autostart = true
